# memogarden.service
[Unit]
Description=MemoGarden Personal Information System
Documentation=https://docs.memogarden.org
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=memogarden
Group=memogarden

# Use CLI wrapper (RFC 004)
ExecStart=/opt/memogarden/bin/memogarden serve

# Working directory
WorkingDirectory=/var/lib/memogarden

# Environment
Environment="PATH=/opt/memogarden/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="MEMOGARDEN_VERB=serve"
Environment="MEMOGARDEN_DATA_DIR=/var/lib/memogarden"
Environment="MEMOGARDEN_CONFIG=/etc/memogarden/config.toml"

# Restart policy
Restart=on-failure
RestartSec=5s
TimeoutSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=memogarden

# Security (RFC 004)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/memogarden /var/log/memogarden /run/memogarden

[Install]
WantedBy=multi-user.target
